# picoCTF-shell-manager/config/shell_manager.target
[Unit]
Description=CTF Shell Manager
Requires=multi-user.target
WantedBy=multi-user.target
